# dumpcap-minimal/DEBIAN/postinst
#!/bin/bash

if ! setcap cap_net_raw,cap_net_admin=ep /usr/bin/dumpcap; then
    echo "Warning: Could not set capabilities on dumpcap."
    echo "You may need to run the following command after booting:"
    echo "sudo setcap cap_net_raw,cap_net_admin=ep /usr/bin/dumpcap"
fi

chown root:wireshark /usr/bin/dumpcap
chmod 750 /usr/bin/dumpcap

systemctl enable dumpcap-setup.service

echo "--------------------------------------------------------"
echo "dumpcap-minimal has been installed"
echo ""
echo "To allow non-root users to capture packets, add them to"
echo "the wireshark group with:"
echo "  sudo usermod -a -G wireshark USERNAME"
echo ""
echo "Then log out and log back in for the changes to take effect."
echo "--------------------------------------------------------"
exit 0
